# systemd/.config/systemd/user/battery-check.service
[Unit]
Description=Battery status checker

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-check.sh
